# deploy/civ0-scheduler.service
[Unit]
Description=Civilization-0 daily scheduler tick (live 0G)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/civilization-0/packages/scheduler
EnvironmentFile=/opt/civilization-0/.env
# pnpm must be on PATH for the service user; adjust if needed.
ExecStart=/usr/bin/env pnpm exec tsx --conditions require scripts/run-scheduler.ts --days 1
