# deploy/civ0-scheduler.timer
[Unit]
Description=Run Civilization-0 scheduler daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
