# deploy/civ0-scheduler.service
[Unit]
Description=Civilization-0 daily scheduler tick (live 0G)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/civilization-0/packages/scheduler
EnvironmentFile=/opt/civilization-0/.env
# Bound a tick: the 0G Storage uploader can wedge in an infinite "Waiting for
# storage node to sync" retry loop when the indexer lags. Without this, a hung
# oneshot stays active forever and systemd never schedules the next fire,
# silently freezing autonomy. On timeout systemd kills the run and the timer
# resumes on its normal cadence (the wedged tick is simply skipped — the
# decision may have been reasoned on 0G but not persisted; no partial DB state).
TimeoutStartSec=600
# pnpm must be on PATH for the service user; adjust if needed.
ExecStart=/usr/bin/env pnpm exec tsx --conditions require scripts/run-scheduler.ts --days 1
# Runs after every tick (success, failure, OR timeout kill) — appends a one-line
# status to /opt/civilization-0/tick.log for at-a-glance autonomy monitoring.
ExecStopPost=/opt/civilization-0/packages/scheduler/scripts/tick-log.sh
